# roles/event_detector/templates/event_detector-supervisor.conf
[program:event_detector]
command=java -jar {{install_dir}}envoy-assembly-0.0.1-SNAPSHOT.jar
directory={{install_dir}}
