# roles/nats_streaming_server/templates/streaming-supervisor.conf
[program:nats-streaming]
command={{executable}}
numprocs=1
autostart=true
autorestart=unexpected
